# ansible/roles/docker_setup/templates/timers/maintenance.service.j2
[Unit]
Description=Runs the maintenance script

[Service]
Type=oneshot
ExecStart={{ docker_script_dir }}/maintenance.sh