# ansible/roles/docker_setup/templates/timers/maintenance.timer.j2
[Unit]
Description=Schedule to run the maintenance script at the night

[Timer]
Persistent=true
OnCalendar={{ item.schedule }}

[Install]
WantedBy=timers.target