# ansible/roles/docker/templates/timers/staleFileHandleHandler.service.j2
[Unit]
Description=Runs the stale file handle handler script

[Service]
Type=oneshot
ExecStart={{ script_dir }}/staleFileHandleHandler.sh