# roles/docker/templates/timers/maintenance.timer.j2
[Unit]
Description=Schedule to run the maintenance script at the night

[Timer]
Persistent=true
OnCalendar=*-*-* 04:00:00

[Install]
WantedBy=timers.target