# roles/docker/templates/timers/maintenance.service.j2
[Unit]
Description=Runs the maintenance script

[Service]
Type=oneshot
ExecStart={{ script_dir }}/maintenance.sh