# ansible/roles/container/docker_setup/templates/timers/staleFileHandleHandler.timer.j2
[Unit]
Description=Schedule to run the stale file handle handler

[Timer]
Persistent=true
OnCalendar={{ item.schedule }}

[Install]
WantedBy=timers.target